# bot.service
[Unit]
Description=Trading Bot Manager (Multi-Pair)
After=network.target

[Service]
Type=simple
# On server (GitHub deploy): User=ubuntu, path /home/ubuntu/trading-bot. Local runs use project root + BOT_HOME from env.
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment=PATH=/home/ubuntu/trading-bot/venv/bin
Environment=BOT_HOME=/home/ubuntu/trading-bot
Environment=BOT_ENTRY_ORDER_MODE=LIMIT_THEN_MARKET
Environment=BOT_ENTRY_LIMIT_TTL_SEC=5
Environment=BOT_POSITION_RECONCILE_SEC=30
# MUST run bot_manager.py (not main.py). The manager spawns one main.py <PAIR> per enabled pair.
ExecStart=/home/ubuntu/trading-bot/venv/bin/python bot/bot_manager.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
