# server.service
[Unit]
Description=Trading Bot Flask API
After=network.target

[Service]
Type=simple
# On server (GitHub deploy): User=ubuntu, path /home/ubuntu/trading-bot.
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment=PATH=/home/ubuntu/trading-bot/venv/bin
Environment=BOT_HOME=/home/ubuntu/trading-bot
ExecStart=/home/ubuntu/trading-bot/venv/bin/python server/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
